# systemd/.config/systemd/user/rclone-dropbox.timer
[Unit]
Description=Run Rclone Dropbox Bisync every 5 minutes

[Timer]
# Trigger 1 minute after the timer unit is started/loaded
OnActiveSec=1min
# Trigger every 5 minutes (at :00, :05, :10, etc.)
OnCalendar=*:0/5
# Ensure it runs if the machine was suspended or off during a scheduled slot
Persistent=true

[Install]
WantedBy=timers.target
